# publish/linux/openprint.service
[Unit]
Description=OpenPrint - Thermal Printer Proxy Service
Documentation=https://github.com/aliexz011/openprint
After=network.target

[Service]
Type=notify
User=openprint
Group=openprint
WorkingDirectory=/opt/openprint
ExecStart=/opt/openprint/OpenPrint
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/openprint

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openprint

[Install]
WantedBy=multi-user.target
